# docker/supervisord.conf
; Runs the whole Aegis stack inside a single container:
;   - the Rust gateway (ingest 5140, REST/MCP/UI 7321)
;   - the self-driving workload microservice (web UI 8080)
; Program output is forwarded to the container's stdout so `docker logs` works.

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:gateway]
command=aegis-daemon --config /app/configs/aegis.docker.toml
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:workload]
command=python -m workload
directory=/app/microservice
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
